# platform/centec-arm64/phytium-bsp/debian/phytium-bsp.postinst
systemctl enable phytium-bsp.service
systemctl start phytium-bsp.service
