# platform/centec-arm64/sonic-platform-modules-phytium/debian/platform-modules-phytium-8t6s10x2c.postinst
systemctl enable 8t6s10x2c_platform.service
systemctl start 8t6s10x2c_platform.service
