# platform/centec-arm64/sonic-platform-modules-g530/debian/platform-modules-g530-sr100-8t4s26x.postinst
systemctl enable sr100_8t4s26x_platform.service
systemctl start sr100_8t4s26x_platform.service
